# var/docker/supervisord/frontend.conf
[program:frontend]
directory=/app
command=pnpm --filter ./apps/frontend run start
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
startsecs=3
